# speedtest/lib/systemd/user/speedtest.timer
[Unit]
Description=speedtest

[Timer]
OnCalendar=*:0:00

[Install]
WantedBy=timers.target